# etc/profile.d/motd.sh
#!/bin/bash
neofetch
ansiweather -l {LOC} -u imperial
